# tools/treble/hacksaw/scripts/hacksaw.socket
[Socket]
ListenStream=/var/run/hacksaw.sock
SocketMode=0660
SocketUser=root
SocketGroup=hacksaw

[Install]
WantedBy=sockets.target
